# roles/node-exporter/files/node-exporter.service
[Unit]
Description=Prometheus Node Exporter

[Service]
User=prometheus
ExecStart=/usr/local/bin/node_exporter

[Install]
WantedBy=default.target
