# roles/node-exporter/files/node-exporter.upstart.conf
# ansible

description     "Prometheus Node Exporter"

start on virtual-filesystems
stop on runlevel [06]

respawn
respawn limit 5 30
limit nofile 65550 65550

# Environment
setuid prometheus
setgid prometheus
console log

exec /usr/local/bin/node_exporter
